# dist/supervisord.ini
[supervisord]
nodaemon=true

[program:rsyslog]
autorestart=true
command=rsyslogd -n
directory=/tmp
stdout_logfile=syslog
stderr_logfile=syslog
user=root

[program:nginx]
autorestart=true
command=nginx -g "daemon off;"
stdout_logfile=syslog
stderr_logfile=syslog
username=http

[program:php-fpm]
command=/usr/sbin/php-fpm7 -F
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

[program:mariadb]
command=mysqld_safe
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
